# DZ_27/provisioning/roles/client/files/backup.timer
[Unit]
Description=Run srservice every 20 second
[Timer]
OnUnitActiveSec=20
Unit=backup.service
[Install]
WantedBy=multi-user.target